# salt/server/apache/_supervisor/supervisor-apache.conf
[program:apache]
command=/usr/sbin/apache2ctl -DFOREGROUND
killasgroup=true
stopasgroup=true
