# salt/website/django/memcached/supervisor-memcached.conf
[program:memcached]
command=/usr/share/memcached/scripts/systemd-memcached-wrapper /etc/memcached.conf
killasgroup=true
stopasgroup=true
